# deploy/one-click/systemd/cube-sandbox-network-agent.service
[Unit]
Description=Cube Sandbox network-agent
After=network-online.target
Wants=network-online.target
ConditionPathExists=/usr/local/services/cubetoolbox/.one-click.env
PartOf=cube-sandbox-control.target cube-sandbox-compute.target

[Service]
Type=simple
User=root
WorkingDirectory=/usr/local/services/cubetoolbox
EnvironmentFile=-/usr/local/services/cubetoolbox/.one-click.env
ExecStart=/usr/bin/bash /usr/local/services/cubetoolbox/scripts/systemd/network-agent-start.sh
ExecStartPost=/usr/bin/bash /usr/local/services/cubetoolbox/scripts/systemd/network-agent-postcheck.sh
Restart=on-failure
RestartSec=2s
# network-agent-postcheck.sh computes a bounded startup wait from tap_init_num.
TimeoutStartSec=infinity
TimeoutStopSec=30s

[Install]
WantedBy=cube-sandbox-control.target
WantedBy=cube-sandbox-compute.target
