# deploy/one-click/systemd/cube-sandbox-webui.service
[Unit]
Description=Cube Sandbox Web UI container
After=docker.service network-online.target cube-sandbox-cube-api.service
Requires=docker.service
Wants=network-online.target cube-sandbox-cube-api.service
ConditionPathExists=/usr/local/services/cubetoolbox/.one-click.env
PartOf=cube-sandbox-control.target

[Service]
Type=simple
User=root
WorkingDirectory=/usr/local/services/cubetoolbox
EnvironmentFile=-/usr/local/services/cubetoolbox/.one-click.env
ExecStartPre=/usr/bin/bash /usr/local/services/cubetoolbox/scripts/systemd/webui-prepare.sh
ExecStart=/usr/bin/bash /usr/local/services/cubetoolbox/scripts/systemd/webui-start.sh
ExecStartPost=/usr/bin/bash /usr/local/services/cubetoolbox/scripts/systemd/webui-postcheck.sh
ExecStop=/usr/bin/bash /usr/local/services/cubetoolbox/scripts/systemd/webui-stop.sh
Restart=on-failure
RestartSec=2s
# Override short distro defaults explicitly.
TimeoutStartSec=60s
TimeoutStopSec=30s

[Install]
WantedBy=cube-sandbox-control.target
