# deploy/one-click/systemd/cube-sandbox-cubemaster.service
[Unit]
Description=Cube Sandbox CubeMaster
After=network-online.target cube-sandbox-mysql.service cube-sandbox-redis.service
Wants=network-online.target cube-sandbox-mysql.service cube-sandbox-redis.service
ConditionPathExists=/usr/local/services/cubetoolbox/.one-click.env
RequiresMountsFor=/data/CubeMaster/storage
PartOf=cube-sandbox-control.target

[Service]
Type=simple
User=root
WorkingDirectory=/usr/local/services/cubetoolbox
EnvironmentFile=-/usr/local/services/cubetoolbox/.one-click.env
Environment=CUBE_MASTER_CONFIG_PATH=/usr/local/services/cubetoolbox/CubeMaster/conf.yaml
ExecStart=/usr/bin/bash /usr/local/services/cubetoolbox/scripts/systemd/cubemaster-start.sh
ExecStartPost=/usr/bin/bash /usr/local/services/cubetoolbox/scripts/systemd/cubemaster-postcheck.sh
Restart=on-failure
RestartSec=2s
# Override short distro defaults explicitly so SIGTERM has room to drain DB
# connections / flush state before SIGKILL.
TimeoutStartSec=60s
TimeoutStopSec=30s

[Install]
WantedBy=cube-sandbox-control.target
